# meta-sigma-racer-wingman/recipes-sigma-racer-wingman/can-network/can-network/sigma-racer-wingman-can-setup.service
[Unit]
Description=Sigma Racer Wingman SocketCAN interface setup
DefaultDependencies=no
Before=sigma-racer-vehicle.service
After=systemd-modules-load.service
Wants=systemd-modules-load.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/sigma-racer-wingman/vehicle.env
ExecStart=/usr/lib/sigma-racer-wingman/can-setup.sh

[Install]
WantedBy=multi-user.target
